# scripts/services/update-notifier.service
[Unit]
Description=Notify available system updates (Pacman + AUR)

[Service]
ExecStart=%h/.local/bin/pacman-update-check.sh


